# deploy/systemd/nginx.service.d/fotoserver.conf
# Drop-in: bindet nginx in den Fotoserver-Lifecycle ein.
# Nach Installation: sudo systemctl daemon-reload
[Unit]
# nginx stoppt/startet mit fotoserver.target
PartOf=fotoserver.target
# Starker Bind: wenn fotoserver.target in failed-State geht, wird nginx mitgestoppt
BindsTo=fotoserver.target
# nginx erst starten, wenn das FastAPI-Backend bereit ist
After=fotoserver-api.service
