# deploy/systemd/hostapd.service.d/fotoserver.conf
# Drop-in: bindet hostapd in den Fotoserver-Lifecycle ein.
# Nach Installation: sudo systemctl daemon-reload
[Unit]
# hostapd stoppt/startet mit fotoserver.target
PartOf=fotoserver.target
# Starker Bind: wenn fotoserver.target in failed-State geht, wird auch der Hotspot gestoppt.
# Kein After=fotoserver.target hier — würde zirkuläre Abhängigkeit erzeugen
# (fotoserver.target After=hostapd.service).
BindsTo=fotoserver.target
