# deploy/systemd/dnsmasq.service.d/fotoserver.conf
# Drop-in: bindet dnsmasq in den Fotoserver-Lifecycle ein.
# Nach Installation: sudo systemctl daemon-reload
[Unit]
# dnsmasq stoppt/startet mit fotoserver.target
PartOf=fotoserver.target
# Starker Bind: wenn fotoserver.target in failed-State geht, wird dnsmasq mitgestoppt
BindsTo=fotoserver.target
# dnsmasq erst nach hostapd starten (Hotspot muss aktiv sein)
After=hostapd.service
